# dev-container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/pulp/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid

[unix_http_server]
file=/var/run/supervisord/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord/supervisor.sock

[program:postgresql]
command=/usr/pgsql-16/bin/postgres -D /var/lib/pgsql/16/data
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/pulp/postgresql-stdout.log
stderr_logfile=/var/log/pulp/postgresql-stderr.log

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --protected-mode yes
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/pulp/redis-stdout.log
stderr_logfile=/var/log/pulp/redis-stderr.log

[program:pulp-api]
command=/usr/bin/pulp-api
autostart=true
autorestart=true
priority=300
environment=PATH="/usr/local/lib/pulp/bin:%(ENV_PATH)s",PULP_SETTINGS="/etc/pulp/settings.py"
stdout_logfile=/var/log/pulp/pulp-api-stdout.log
stderr_logfile=/var/log/pulp/pulp-api-stderr.log

[program:pulp-content]
command=/usr/bin/pulp-content
autostart=true
autorestart=true
priority=300
environment=PATH="/usr/local/lib/pulp/bin:%(ENV_PATH)s",PULP_SETTINGS="/etc/pulp/settings.py"
stdout_logfile=/var/log/pulp/pulp-content-stdout.log
stderr_logfile=/var/log/pulp/pulp-content-stderr.log

[program:pulp-worker]
command=/usr/bin/pulp-worker
autostart=true
autorestart=true
priority=300
environment=PATH="/usr/local/lib/pulp/bin:%(ENV_PATH)s",PULP_SETTINGS="/etc/pulp/settings.py"
stdout_logfile=/var/log/pulp/pulp-worker-stdout.log
stderr_logfile=/var/log/pulp/pulp-worker-stderr.log
